# Scripts/crontab
# To be added to crontab

@reboot sh /home/pi/tvshutdown.sh 2>&1


